# services/user/xbanish.service
[Unit]
Description=Run ksuperkey a daemon

[Service]
Type=simple
ExecStart=/usr/bin/xbanish

[Install]
WantedBy=default.target
